# script/mihomo.service
[Unit]
Description=clash auto run
 
[Service]
Type=simple
 
ExecStart=/usr/local/bin/mihomo -d /etc/mihomo/
 
[Install]
WantedBy=default.target